# files/volumepad/10-volumepad.rules
KERNEL=="event*", SUBSYSTEM=="input", ATTRS{name}=="Dean Camera LUFA Keyboard Demo", RUN+="/usr/local/bin/volumepad-control %k"
